# shutdown-button.service
[Unit]
Description=Shutdown Button Service

[Service]
Type=simple
ExecStart=/usr/local/bin/shutdown_button.py

[Install]
WantedBy=default.target
